# 7.2/aws/bin/init.d/jbosseap7.service
[Unit]
Description=JBoss EAP Systemctl script
After=NetworkManager.service

[Service]
Type=simple
User=ec2-user
ExecStart=/bin/sh /home/ec2-user/jboss-eap-7.2/bin/init.d/ec2-standalone.sh
Restart=always
RestartSec=5
TimeoutStartSec=5
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target